# profile.d/editor.sh
#!/bin/sh

export EDITOR="emacs -nw"
export ALTERNATE_EDITOR="ee"
