# profile.d/go.sh
#!/bin/sh

export GOPATH=$HOME/gopath
export GOROOT=$HOME/go

export PATH=$PATH:$GOROOT/bin:$GOPATH/bin
