# profile.d/nash.sh
#!/bin/sh

NASH=$HOME/nash
NASHROOT=$HOME/nashroot

export PATH=$PATH:$NASHROOT/bin:$NASH/bin
