# profile.d/nash.sh
#!/bin/sh

NASH=$HOME/nashpath
NASHROOT=$HOME/nash

export PATH=$PATH:$NASHROOT/bin:$NASH/bin
